# .config/production/supervisord.conf
[program:nginx]
command=nginx

[program:uwsgi]
command=uwsgi --ini /srv/project/.config/production/uwsgi.ini
